# debian/postinst
#!/bin/bash

systemctl daemon-reload
systemctl unmask galliumos-skylake
systemctl enable galliumos-skylake
systemctl start galliumos-skylake

exit 0
